# template/etc/systemd/system/redis-server.service.d/zzz-zzz-fix-perms.conf
[Service]
ReadWritePaths=-/etc/redis
ReadWritePaths=-/var/lib/redis
ReadWritePaths=-/var/log/redis
ReadWritePaths=-/var/run/redis
